# root/etc/services.d/wireguard/run
#!/usr/bin/with-contenv bash

_term() {
  echo "Caught SIGTERM signal!"
  wg-quick down wg1
}

trap _term SIGTERM

wg-quick up wg1

sleep infinity &

wait
